# examples/crontab
*/5 * * * * cd /path/to/all_tweets_must_die && rake run &> /tmp/all_tweets_must_die.log